# misc/mba-6-1/fand.service
[Unit]
Description=applesmc fan controller

[Service]
Type=simple
ExecStart=/usr/sbin/fand

[Install]
WantedBy=multi-user.target
