# misc/zram/enable-zram@.service
[Unit]
Description=enable zram

[Service]
Type=oneshot
ExecStart=/usr/sbin/enable-zram %I

[Install]
WantedBy=multi-user.target
