# misc/oneplus6/80-wlan0-mac-address.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/usr/bin/ip link set dev $name address 96:2a:56:07:0b:9c"
